# dist/corepost-agent.service
[Unit]
Description=CorePost Agent
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/corepost-agent.env
ExecStart=/usr/local/lib/corepost-agent/corepost-agent.sh
Restart=always
RestartSec=3

# Hardening (keep minimal: agent may need system tools like loginctl/systemctl).
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
